# deploy/systemd/osint-agent-network-backup.service
[Unit]
Description=OSINT Agent Network Backup

[Service]
Type=oneshot
WorkingDirectory=/home/aidi/apps/osint-agent-network
EnvironmentFile=-/home/aidi/apps/osint-agent-network/.env
ExecStart=/bin/bash /home/aidi/apps/osint-agent-network/scripts/backup.sh
