# deploy/systemd/osint-agent-network-backup.timer
[Unit]
Description=Daily OSINT Agent Network Backup

[Timer]
OnCalendar=*-*-* 03:30:00
Persistent=true
Unit=osint-agent-network-backup.service

[Install]
WantedBy=timers.target
