# deploy/systemd/osint-agent-network-api.service
[Unit]
Description=OSINT Agent Network API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/osint-agent-network
EnvironmentFile=/opt/osint-agent-network/.env
Environment=PYTHONPATH=/opt/osint-agent-network/backend
ExecStart=/usr/bin/python3 -m app.main
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
